# supervisord.conf
[supervisord]
nodaemon=true # false  supervisord -c /etc/supervisord.conf (启动方式)

[program:nginx]
command=/usr/sbin/nginx